# systemd/solbot-selftest.timer
[Unit]
Description=Run solbot self-test daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
